# deploy/linux/agopenweb.service
# AgOpenWeb headless host — systemd service unit (Phase 10 appliance deployment).
#
# Runs the display-less .NET host as a daemon on a cab-PC / one-box SBC. The host
# has NO native UI; operators use a browser (remote tablet, or a local
# chromium-kiosk pointed at http://localhost:5174). See Plans/DEPLOYMENT_PATTERNS.md.
#
# Install with deploy/linux/install.sh, which publishes the self-contained build to
# /opt/agopenweb, creates the agopenweb service user, and installs + enables this unit.
# After editing this file by hand: `sudo systemctl daemon-reload && sudo systemctl restart agopenweb`.

[Unit]
Description=AgOpenWeb headless GPS guidance host (browser UI on :5174)
Documentation=https://github.com/chriskinal/AgOpenWeb
# The host opens UDP sockets to the AiO / GPS modules and (optionally) an NTRIP
# caster, so wait for the network to be actually up, not just configured.
After=network-online.target
Wants=network-online.target

[Service]
# Type=notify: the .NET host (UseSystemd) sends READY=1 once it is serving, so
# systemd only considers the unit "started" when the browser endpoint is live —
# and the WatchdogSec timer below only starts counting after READY.
Type=notify
NotifyAccess=main

ExecStart=/opt/agopenweb/AgOpenWeb.Desktop
WorkingDirectory=/opt/agopenweb

# Dedicated unprivileged service account. Device-group membership (dialout for USB
# serial GPS/AiO; can where a SocketCAN group exists) is granted to the user by
# install.sh via usermod, and systemd applies it automatically at spawn from the
# group database (initgroups) because User= is set. We deliberately do NOT list
# SupplementaryGroups= here: systemd hard-fails the unit (status=216/GROUP) if any
# named group is missing, and the SocketCAN "can" group does not exist on stock
# Debian (incl. the Uno Q) — CAN frame I/O needs no unix group anyway.
User=agopenweb
Group=agopenweb

# Field data + config live OUTSIDE the program dir so a program update (which replaces
# /opt/agopenweb) never wipes the operator's fields/tools/vehicles/config. They live in
# /home/agopenweb — the agopenweb user's home, where people look for user data (it sits in
# /home alongside the real users). install.sh creates it (0755, owned by agopenweb).
# UMask=0022 keeps everything the daemon writes WORLD-READABLE but daemon-write-only: any
# login user can browse/back up the fields, only the service can modify them. AGOPENWEB_DATA
# is the explicit data-root the app honors (see AppDataRoot) → /home/agopenweb/AgOpenWeb.
UMask=0022
Environment=AGOPENWEB_DATA=/home/agopenweb
Environment=HOME=/home/agopenweb
# .NET diagnostics pipe is pointless on an appliance and trips some hardening.
Environment=DOTNET_EnableDiagnostics=0

# Keep the brain alive. RestartSec gives modules a moment to settle between crashes.
Restart=always
RestartSec=2
# Watchdog: the host pets it via sd_notify at half this interval
# (SystemdWatchdogService). A wedged host misses the deadline → systemd kills and
# restarts it. Raise if you see false trips under heavy field load.
WatchdogSec=30
TimeoutStartSec=60
TimeoutStopSec=15

# Logs → journald (view with: journalctl -u agopenweb -f).
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agopenweb

# --- Hardening. Loosened only where the workload needs it; tighten per site. ---
NoNewPrivileges=true
ProtectSystem=full
# Must be false: the data home is /home/agopenweb, which ProtectHome would hide/lock.
ProtectHome=false
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
# The .NET layer talks to modules over UDP (AF_INET) and to NTRIP over TCP; AF_UNIX
# is used internally; AF_CAN is listed so a future SocketCAN path is not blocked.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_CAN

[Install]
WantedBy=multi-user.target
